# docs/systemd/dhcp-log.service
[Unit]
Description=Write dhcp events to Rstat
After=syslog.target mariadb.service dnsmasq.service dhcpd.service
PartOf=dnsmasq.service

[Service]
Type=forking
ExecStart=/usr/local/scripts/dhcp-log.pl --start
Restart=always
PIDFile=/var/run/dhcp-log.pid

[Install]
WantedBy=multi-user.target
