# docs/systemd/syslog-stat.service
[Unit]
Description=Write syslog to Eye
After=network.target mariadb.service syslog-ng.service
PartOf=syslog-ng.service

[Service]
Type=simple
User=eye
Group=eye
ExecStart=/opt/Eye/scripts/syslog-stat.pl
RuntimeDirectory=eye-syslog
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
RestartSec=30
Restart=on-failure

[Install]
WantedBy=multi-user.target
